# supervisord/sig2mqtt.conf
[program:sig2mqtt]
command=node /home/pi/prj/sig2mqtt/sig2mqtt.js
directory=/home/pi/
user=pi
autostart=true
autorestart=false
